# linux/systemd/claude-token-refresh.service
[Unit]
Description=Refresh Claude OAuth Token
After=network-online.target
Wants=network-online.target
# Don't fail if network not ready - script handles retries

[Service]
Type=oneshot
ExecStart=~/.claude/scripts/refresh-claude-token.sh
Environment="CLAUDE_HOME=~/.claude"
# Note: For user-level service, HOME is set automatically
# For system-level, override with: Environment="HOME=/home/username"
StandardOutput=journal
StandardError=journal
# Prevent running too frequently on repeated failures
RestartSec=60
